# stage3/03-wazigate/files/reconnect_wifi.timer
[Unit]
Description=Reconnect wifi service timer

[Timer]
OnUnitActiveSec=120s
OnBootSec=120s

[Install]
WantedBy=timers.target